# rhodes.service
[Unit]
Description=RHODES — Agentic Infrastructure Operations
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/path/to/rhodes
Environment=PATH=/path/to/node/bin:/usr/bin:/bin
EnvironmentFile=/path/to/rhodes/.env
ExecStart=/path/to/node/bin/npx tsx src/index.ts full
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
